# stage4/04-custom-mods/files/xdg/autostart/noscreenblank.desktop
[Desktop Entry]
Type=Application
Exec=xset s off
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
Name[en_US]=noscreenblank
Name=noscreenblank
Comment[en_US]=
Comment=
